# ansible/roles/shlink/templates/shlink/shlink.sh.j2
#!/bin/sh

# PROVIDE: shlink
# REQUIRE: LOGIN network
# BEFORE:  DAEMON
# KEYWORD: shutdown

{% for line in warning_comment %}
# {{ line }}
{% endfor %}

. /etc/rc.subr

name="shlink"
rcvar="shlink_enable"
pidfile="/var/run/${name}.pid"

load_rc_config $name

: ${shlink_enable:="NO"}
: ${shlink_logfile:="/var/log/shlink.log"}
: ${shlink_subprocess:="/usr/local/shlink/bin/rr"}
: ${shlink_subprocess_flags:="serve -c /usr/local/shlink/config/roadrunner/.rr.yml"}

command="/usr/sbin/daemon"
command_args="-f -p ${pidfile} -o ${shlink_logfile} -- ${shlink_subprocess} ${shlink_subprocess_flags}"
procname="${shlink_subprocess}"

start_precmd="shlink_prestart"

shlink_prestart() {
    cd /usr/local/shlink || exit 1
    return 0
}

# Default settings
: ${shlink_enable:="NO"}

run_rc_command "$1"
