# ansible/roles/expense-analysis-app/templates/expense-analysis.sh.j2
#!/bin/sh

# PROVIDE: expense_analysis
# REQUIRE: LOGIN network
# BEFORE:  DAEMON
# KEYWORD: shutdown

{% for line in warning_comment %}
# {{ line }}
{% endfor %}

. /etc/rc.subr

name="expense_analysis"
rcvar="expense_analysis_enable"
pidfile="/var/run/${name}.pid"

load_rc_config $name

: ${expense_analysis_enable:="NO"}
: ${expense_analysis_logfile:="/var/log/expense-analysis/server.log"}
: ${expense_analysis_dir:="{{ expense_analysis_dir }}"}

_expense_analysis_subprocess_user="{{ expense_analysis_user }}"
_expense_analysis_subprocess="poetry"
_expense_analysis_subprocess_flags="run python -m expense_analysis.app"
_expense_analysis_subprocess_env_file="{{ expense_analysis_env_file }}"

command="/usr/sbin/daemon"
command_args="-f -p ${pidfile} -o ${expense_analysis_logfile} -u ${_expense_analysis_subprocess_user} -- ${_expense_analysis_subprocess} ${_expense_analysis_subprocess_flags}"
procname="${expense_analysis_dir}/backend/.venv/bin/python"

start_precmd="expense_analysis_prestart"

expense_analysis_prestart() {
    cd ${expense_analysis_dir}/backend || exit 1

    # Ensure log directory exists
    mkdir -p "$(dirname "${expense_analysis_logfile}")"
    chmod u+rw,go+r "${expense_analysis_logfile}"
    return 0
}

export ENV_FILE="${_expense_analysis_subprocess_env_file}"

run_rc_command "$1"
